# contrib/ds-watch.service
[Unit]
Description=ds-watch: daily CZDS DS collector run
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/Dokumente/03_Projekte/03_IETF/01_EvDANE/ds-watch
ExecStart=%h/Dokumente/03_Projekte/03_IETF/01_EvDANE/ds-watch/.venv/bin/ds-watch run
# Exit code 2 (quarantine / grant expired / new ToU) deliberately stays a
# visible failure: `systemctl --user status ds-watch` shows it.
