# contrib/ds-watch.timer
[Unit]
Description=ds-watch daily after the CZDS regeneration window (00:00-06:00 UTC)

[Timer]
OnCalendar=*-*-* 06:30:00 UTC
RandomizedDelaySec=15m
# Catch up on missed runs (laptop off/suspended) at next boot;
# the HEAD/24h guard in the client prevents duplicate downloads.
Persistent=true

[Install]
WantedBy=timers.target
